# data/locationd.conf.in
description "Location Services"

start on started dbus

respawn

exec /usr/bin/locationd run
